# files/etc/supervisor.d/pgadmin.ini
[program:pgadmin]
command=/usr/bin/uwsgi --ini /etc/supervisor.d/pgadmin.cfg
priority=100
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
